# service_files/replicant.timer
[Unit]
Description=Run the latest release script every 60 minutes

[Timer]
OnBootSec=10min
OnUnitActiveSec=60min

[Install]
WantedBy=timers.target
